# zeroconf-avahi/llmnr.service
[Unit]
Description=LLMNR
Requires=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/llmnr.py
TimeoutSec=infinity

[Install]
WantedBy=multi-user.target
